# package/base-files/files/etc/crontabs/root
*/10 * * * * /etc/init.d/sysntpd restart
